# systemd/user/dunst.service
[Unit]
Description=Dunst

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=wm.target
