# systemd/user/polybar.service
[Unit]
Description=Polybar

[Service]
ExecStart=/usr/bin/polybar

[Install]
WantedBy=wm.target
