# systemd/user/wm.target
[Unit]
Description=Window manager target
